# recipes-extended/ostree/files/ostree-pending-reboot.path
[Unit]
Description=Aktualizr OSTree Update Automatic Reboot

[Path]
PathExists=/run/need-reboot

[Install]
WantedBy=multi-user.target
