# recipes-extended/ostree/files/ostree-pending-reboot.service
[Unit]
Description=Aktualizr OSTree Update Automatic Reboot

[Service]
Type=simple
ExecStart=/usr/sbin/reboot
