# recipes-sota/ostree/files/ostree-pending-reboot.path
[Unit]
Description=Aktualizr OSTree Update Automatic Reboot

[Path]
PathExists=/run/aktualizr-session/need_reboot

[Install]
WantedBy=multi-user.target
